# docker/httpd.service
.include /lib/systemd/system/httpd.service
[Service]
UMask=002
ExecStartPre=/usr/local/bin/install-wordpress.sh
